# scripts/templates/freepbx-module-upgrade.timer
[Unit]
Description=Actualizacion semanal de modulos FreePBX (domingo 03:00)

[Timer]
OnCalendar=Sun *-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
