# system_files/usr/lib/systemd/system/bootc-fetch-apply-updates.timer.d/10-schedule.conf
# Update del SO + reinicio una sola vez al dia, a las 04:00 hora local.
# "bootc update --apply" NO hace nada (ni reinicia) si no hay imagen nueva.
# Se vacian primero los disparadores de fabrica para que el UNICO sea las 04:00.
[Timer]
OnBootSec=
OnUnitInactiveSec=
OnUnitActiveSec=
# Ojo: los disparadores tipo lista se limpian con asignacion vacia, pero este
# retardo escalar exige el 0 explicito (si no, hereda las 2h de fabrica).
RandomizedDelaySec=0
OnCalendar=*-*-* 04:00:00
Persistent=true
