# system_files/usr/lib/systemd/system/systemd-nspawn@freepbx.service.d/10-limits.conf
# Tope de memoria del contenedor FreePBX (equipo de 4GB compartido con Kanidm y el host).
# MemoryHigh frena suave (reclaim); MemoryMax es el limite duro.
[Service]
MemoryHigh=1400M
MemoryMax=1800M
