# system_files/usr/lib/systemd/system/pbx-idm-backup.timer
[Unit]
Description=Backup nocturno consolidado (02:30)

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
