# system_files/usr/lib/systemd/system/tailscale-cert.timer
[Unit]
Description=Renovacion diaria del certificado TLS de Tailscale (03:30)

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
